# puppet/files/restart-browser.service
[Unit]

[Service]
ExecStart=/usr/local/sbin/restart-browser.sh

[Install]
WantedBy=default.target
